# cakephp/gearmanworker/qlts.ini
[program:export-order]
process_name=%(program_name)s_%(process_num)02d
command=php -f /var/www/html/asset/app/Console/cake.php gearman_worker main export
autostart=true
autorestart=true
numprocs=8
